# deploy/million-item-selector.service
[Unit]
Description=Million Item Selector Express application
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/html/test
Environment=NODE_ENV=production
Environment=HOST=127.0.0.1
Environment=PORT=4173
ExecStart=/usr/bin/node /var/html/test/src/server.js
Restart=on-failure
RestartSec=3
TimeoutStopSec=15
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
